# sysdrv/rootfs_overlay/lib/udev/rules.d/99-sd-test.rules
KERNEL=="mmcblk1p*[0-9]", ACTION=="add", RUN+="/bin/sh /etc/auto_test.sh add "
KERNEL=="mmcblk1p*[0-9]", ACTION=="remove", RUN+="/bin/sh /etc/auto_test.sh remove "
